# spamhaus_drop/spamhaus_drop.cron
45 03 * * * root /etc/nftfw/spamhaus_drop.d/getdrop
