# debian/preinst
#!/bin/sh
# preinst script for nftfw
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    install)
	# for an initial install need to check if a previous manual
	# installation for systemd is running. Cancel installation
	# until it goes away.
	if [ ! -n "$2" ]; then
	    servicefile=$(systemctl show -p FragmentPath nftfw.service | sed -e 's/^.*=//') || true
	    if [ "$file" != "" -a -f "$file" ]; then
		check=$(grep '/usr/local/bin' $servicefile) || true
		if [ "$check" != '' ]; then
		    pathfile=$(systemctl show -p FragmentPath nftfw.path | sed -e 's/^.*=//') || true
		    cat >&2 <<EOF
The systemd daemon is running a manually installed version of nftfw.service and
nftfw.path. This must be removed before this package can be installed.
Run these commands:
   sudo systemctl stop nftfw.path nftfw.service
   sudo systemctl disable nftfw.path nftfw.service
   sudo rm $servicefile $pathfile
   sudo systemctl daemon reload
EOF
		    exit 1
		fi
	    fi
	fi
    ;;

    upgrade)
    ;;
    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
